# watch_directory.service
[Unit]
Description=Watch Directory Service

[Service]
Type=simple
ExecStart= /usr/local/bin/vladlen_script_10.sh

[Install]
WantedBy=default.target
